# templates/default/nodejs.systemd.service.erb
[Service]
ExecStart=/usr/bin/env node <%= @entry %>
User=<%= @user %>
Group=<%= @group %>
